# dirflow.service
[Unit]
Description=DirFlow: Lightweight Local File Organizer Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/your/dirflow
ExecStart=/usr/bin/python organizer.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
